# .config/systemd/user/offlineimap.timer
[Timer]
OnUnitInactiveSec=5min
Unit=offlineimap.service

[Install]
WantedBy=default.target
